# Roles/restic/templates/systemd/restic-scan@.service.j2
[Unit]
Description=Restic Snapshot Scan Service for %i
After=network-online.target
Wants=network-online.target
Documentation=https://restic.readthedocs.io/

[Service]
Type=oneshot
User=root
Group=root

# Environment
Environment="RESTIC_REPOSITORY={{ restic_repository }}"
Environment="RESTIC_PASSWORD_FILE={{ restic_config_dir }}/passwords/{{ restic_password_file }}"
{% if restic_backend_type == "s3" %}
Environment="AWS_ACCESS_KEY_ID={{ restic_s3_access_key }}"
Environment="AWS_SECRET_ACCESS_KEY={{ restic_s3_secret_key }}"
Environment="AWS_DEFAULT_REGION={{ restic_s3_region }}"
{% endif %}

# CheckMK Custom Tag
Environment="CHECKMK_TAG=restic_scan_%i"

# Intelligent lock check - only unlock stale locks, let --retry-lock handle fresh locks
{% set lock_max_age = item.lock_max_age_hours | default(restic_lock_max_age_hours) %}
ExecStartPre=/bin/bash -c ' \
  MAX_AGE_SECONDS=$(({{ lock_max_age }} * 3600)); \
  LOCKS=$(restic list locks --no-lock 2>/dev/null || echo ""); \
  \
  if [ -n "$LOCKS" ]; then \
    LOCK_TIME=$(echo "$LOCKS" | head -1 | sed -n "s/.*at \([0-9]\{4\}-[0-9]\{2\}-[0-9]\{2\} [0-9]\{2\}:[0-9]\{2\}:[0-9]\{2\}\).*/\1/p" | head -1); \
    if [ -n "$LOCK_TIME" ]; then \
      LOCK_TS=$(date -d "$LOCK_TIME" +%%s 2>/dev/null || echo 0); \
      if [ $LOCK_TS -gt 0 ]; then \
        CURRENT_TS=$(date +%%s); \
        LOCK_AGE=$((CURRENT_TS - LOCK_TS)); \
        \
        if [ $LOCK_AGE -gt $MAX_AGE_SECONDS ]; then \
          logger -t restic-scan -p user.warning "Removing stale lock for %i (age: ${LOCK_AGE}s > ${MAX_AGE_SECONDS}s)"; \
          restic unlock || true; \
        else \
          logger -t restic-scan -p user.info "Lock for %i exists (age: ${LOCK_AGE}s), will use --retry-lock to wait"; \
        fi; \
      fi; \
    fi; \
  fi'

# Execute snapshot scan and save stats
ExecStart=/bin/bash {{ restic_scripts_dir }}/restic-scan-%i.sh

{% if restic_enable_checkmk | default(true) | bool %}
# Send CheckMK status
ExecStopPost=/bin/bash {{ restic_scripts_dir }}/restic-checkmk-status.sh %i scan $EXIT_STATUS
{% endif %}

# Syslog
ExecStopPost=/bin/bash -c ' \
  if [ "$EXIT_STATUS" = "0" ]; then \
    logger -t restic-scan -p user.info "Snapshot scan completed for %i"; \
  else \
    logger -t restic-scan -p user.err "Snapshot scan failed for %i (exit: $EXIT_STATUS)"; \
  fi'

# Resource limits
CPUQuota={{ restic_cpu_quota }}%
IOWeight={{ restic_io_weight }}
Nice={{ restic_nice_level }}

# Security
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths={{ restic_config_dir }} {% if restic_enable_checkmk | default(true) | bool %}{{ restic_checkmk_spool_dir }}{% endif %}

[Install]
WantedBy=multi-user.target
