# Roles/restic/templates/systemd/restic-prune@.service.j2
[Unit]
Description=Restic Prune Service for %i
After=network-online.target
Wants=network-online.target
Documentation=https://restic.readthedocs.io/

[Service]
Type=oneshot
User=root
Group=root

# Environment
Environment="RESTIC_REPOSITORY={{ restic_repository }}"
Environment="RESTIC_PASSWORD_FILE={{ restic_config_dir }}/passwords/{{ restic_playbook_password_file }}"
{% if restic_backend_type == "s3" %}
Environment="AWS_ACCESS_KEY_ID={{ restic_s3_access_key }}"
Environment="AWS_SECRET_ACCESS_KEY={{ restic_s3_secret_key }}"
Environment="AWS_DEFAULT_REGION={{ restic_s3_region }}"
{% endif %}

# CheckMK Custom Tag
Environment="CHECKMK_TAG=restic_prune_%i"

# Restic native lock cleanup (remove stale repository locks)
# Restic's unlock command removes locks older than 30 minutes
ExecStartPre=/bin/bash -c 'restic unlock || true'

# Execute prune with retention policy
ExecStart=/usr/bin/restic forget \
  --keep-last {{ restic_retention_keep_last }} \
  --keep-hourly {{ restic_retention_keep_hourly }} \
  --keep-daily {{ restic_retention_keep_daily }} \
  --keep-weekly {{ restic_retention_keep_weekly }} \
  --keep-monthly {{ restic_retention_keep_monthly }} \
  --keep-yearly {{ restic_retention_keep_yearly }} \
  --prune

# Send CheckMK status
ExecStopPost=/bin/bash {{ restic_scripts_dir }}/checkmk-status.sh %i prune $EXIT_STATUS

# Syslog
ExecStopPost=/bin/bash -c ' \
  if [ "$EXIT_STATUS" = "0" ]; then \
    logger -t restic-prune -p user.info "Prune completed successfully for %i"; \
  else \
    logger -t restic-prune -p user.err "Prune failed for %i (exit: $EXIT_STATUS)"; \
  fi'

# Resource limits
CPUQuota={{ restic_cpu_quota }}%
IOWeight={{ restic_io_weight }}
Nice={{ restic_nice_level }}

# Security
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths={{ restic_config_dir }} {{ restic_checkmk_spool_dir }}

[Install]
WantedBy=multi-user.target
