# Roles/restic/templates/systemd/restic-backup@.service.j2
[Unit]
Description=Restic Backup Service for %i
After=network-online.target
Wants=network-online.target
Documentation=https://restic.readthedocs.io/

[Service]
Type=oneshot
User=root
Group=root

# Environment
Environment="RESTIC_REPOSITORY={{ restic_repository }}"
Environment="RESTIC_PASSWORD_FILE={{ restic_config_dir }}/passwords/{{ restic_password_file }}"
{% if restic_backend_type == "s3" %}
Environment="AWS_ACCESS_KEY_ID={{ restic_s3_access_key }}"
Environment="AWS_SECRET_ACCESS_KEY={{ restic_s3_secret_key }}"
Environment="AWS_DEFAULT_REGION={{ restic_s3_region }}"
{% endif %}
Environment="RESTIC_READ_CONCURRENCY={{ restic_read_concurrency }}"
Environment="RESTIC_LIMIT_UPLOAD={{ restic_upload_limit_kbps }}"
Environment="RESTIC_LIMIT_DOWNLOAD={{ restic_download_limit_kbps }}"

# Restic Hooks (native hook system)
# Restic automatically calls these scripts before/after backup
{% set hook_dir = item.hook_script_dir | default(restic_hooks_dir) %}
Environment="RESTIC_PRE_SCRIPT={{ hook_dir }}/pre-backup-%i.sh"
Environment="RESTIC_POST_SCRIPT={{ hook_dir }}/post-backup-%i.sh"

# CheckMK Custom Tag
Environment="CHECKMK_TAG=restic_backup_%i"

# Intelligent lock check - only unlock stale locks, let --retry-lock handle fresh locks
{% set lock_max_age = item.lock_max_age_hours | default(restic_lock_max_age_hours) %}
ExecStartPre=/bin/bash -c ' \
  MAX_AGE_SECONDS=$(({{ lock_max_age }} * 3600)); \
  LOCKS=$(restic list locks --no-lock 2>/dev/null || echo ""); \
  \
  if [ -n "$LOCKS" ]; then \
    LOCK_TIME=$(echo "$LOCKS" | head -1 | grep -oP "(?<=at )[0-9]{4}-[0-9]{2}-[0-9]{2} [0-9]{2}:[0-9]{2}:[0-9]{2}" | head -1); \
    if [ -n "$LOCK_TIME" ]; then \
      LOCK_TS=$(date -d "$LOCK_TIME" +%%s 2>/dev/null || echo 0); \
      if [ $LOCK_TS -gt 0 ]; then \
        CURRENT_TS=$(date +%%s); \
        LOCK_AGE=$((CURRENT_TS - LOCK_TS)); \
        \
        if [ $LOCK_AGE -gt $MAX_AGE_SECONDS ]; then \
          logger -t restic-backup -p user.warning "Removing stale lock for %i (age: ${LOCK_AGE}s > ${MAX_AGE_SECONDS}s)"; \
          restic unlock || true; \
        else \
          logger -t restic-backup -p user.info "Lock for %i exists (age: ${LOCK_AGE}s), will use --retry-lock to wait"; \
        fi; \
      fi; \
    fi; \
  fi'

# Execute backup (restic will call pre/post hooks automatically)
TimeoutStartSec={{ item.timeout_seconds | default(restic_backup_timeout_seconds) }}
ExecStart=/bin/bash {{ restic_scripts_dir }}/backup-%i.sh

# Send CheckMK status
ExecStopPost=/bin/bash {{ restic_scripts_dir }}/checkmk-status.sh %i backup $EXIT_STATUS

# Syslog on failure
ExecStopPost=/bin/bash -c 'if [ "$EXIT_STATUS" != "0" ]; then logger -t restic-backup -p user.err "Backup failed for %i (exit: $EXIT_STATUS)"; fi'

# Resource limits
CPUQuota={{ restic_cpu_quota }}%
IOWeight={{ restic_io_weight }}
Nice={{ restic_nice_level }}

# Security
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths={{ restic_config_dir }} {{ restic_scripts_dir }} {{ restic_hooks_dir }} {{ restic_checkmk_spool_dir }} {% for source in restic_backup_sources %}{{ source.path }} {% endfor %}

# Restart behavior - retry on failure (e.g., locked repository)
Restart=on-failure
RestartSec={{ restic_restart_sec }}
{% if restic_restart_max_attempts > 0 %}
StartLimitBurst={{ restic_restart_max_attempts }}
StartLimitIntervalSec=0
{% endif %}

[Install]
WantedBy=multi-user.target
