# Roles/restic/templates/systemd/restic-check@.service.j2
[Unit]
Description=Restic Repository Check Service for %i
After=network-online.target
Wants=network-online.target
Documentation=https://restic.readthedocs.io/

[Service]
Type=oneshot
User=root
Group=root

# Environment
Environment="RESTIC_REPOSITORY={{ restic_repository }}"
Environment="RESTIC_PASSWORD_FILE={{ restic_config_dir }}/passwords/{{ restic_playbook_password_file }}"
{% if restic_backend_type == "s3" %}
Environment="AWS_ACCESS_KEY_ID={{ restic_s3_access_key }}"
Environment="AWS_SECRET_ACCESS_KEY={{ restic_s3_secret_key }}"
Environment="AWS_DEFAULT_REGION={{ restic_s3_region }}"
{% endif %}

# CheckMK Custom Tag
Environment="CHECKMK_TAG=restic_check_%i"

# Lock file settings
Environment="RESTIC_LOCK_FILE={{ restic_lock_dir }}/check-%i.lock"
Environment="RESTIC_LOCK_MAX_AGE_HOURS={{ restic_lock_max_age_hours }}"

# Unlock mechanism
ExecStartPre=/bin/bash -c ' \
  if [ -f "${RESTIC_LOCK_FILE}" ]; then \
    LOCK_AGE=$(($(date +%%s) - $(stat -c %%Y "${RESTIC_LOCK_FILE}"))); \
    MAX_AGE=$((${RESTIC_LOCK_MAX_AGE_HOURS} * 3600)); \
    if [ ${LOCK_AGE} -gt ${MAX_AGE} ]; then \
      logger -t restic-check -p user.warning "Stale lock detected for %i, removing"; \
      rm -f "${RESTIC_LOCK_FILE}"; \
      restic unlock || true; \
    else \
      logger -t restic-check -p user.error "Lock file exists for %i, aborting"; \
      exit 1; \
    fi \
  fi'

# Create lock file
ExecStartPre=/bin/bash -c 'touch "${RESTIC_LOCK_FILE}"'

# Execute repository check
ExecStart=/usr/bin/restic check {% if restic_check_read_data %}--read-data{% endif %}

# Remove lock file
ExecStopPost=/bin/bash -c 'rm -f "${RESTIC_LOCK_FILE}"'

# Send CheckMK status
ExecStopPost=/bin/bash {{ restic_scripts_dir }}/checkmk-status.sh %i check $EXIT_STATUS

# Syslog
ExecStopPost=/bin/bash -c ' \
  if [ "$EXIT_STATUS" = "0" ]; then \
    logger -t restic-check -p user.info "Repository check passed for %i"; \
  else \
    logger -t restic-check -p user.err "Repository check failed for %i (exit: $EXIT_STATUS)"; \
  fi'

# Resource limits
CPUQuota={{ restic_cpu_quota }}%
IOWeight={{ restic_io_weight }}
Nice={{ restic_nice_level }}

# Security
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths={{ restic_config_dir }} {{ restic_lock_dir }} {{ restic_checkmk_spool_dir }}

[Install]
WantedBy=multi-user.target
